# DEBIAN/postinst
#!/bin/sh
set -e

systemctl daemon-reload
systemctl enable raid_daemon
systemctl start raid_daemon

